# worker/convertVideo1.conf
[program:convertVideo]
command=php /var/www/html/transcoder/yii rabbitmq/consume convertVideo1 --messagesLimit=250 --memoryLimit=32 --debug --withoutSignals
process_name=%(program_name)s%(process_num)02d
numprocs=1
startsecs=2
autorestart=true
stopsignal=INT
stopasgroup=true
stopwaitsecs=60
stdout_logfile=/var/www/html/transcoder/runtime/logs/rabbitmq.log
stderr_logfile=/var/www/html/transcoder/runtime/logs/rabbitmq.log
